# buildroot-2022.11.1/ajoutprojet/overlay/var/spool/cron/crontabs/root
0 */1 * * * ./usr/bin/partie_velib/database.py
