# buildroot-2022.11.1/ajoutprojet/overlay/var/spool/cron/crontabs/root
0 */2 * * * ./home/satrana/database.py 
